# systemd/nonsense.slice
[Unit]
Description=Nonsense namespace engine top level slice
